# pick-files-daily.service
[Unit]
Description = Copy files every day

[Service]
ExecStart = /snap/bin/pick-files --config ${HOME}/pick-files-config.yaml --verbose
